# unconfinehomeusers.path
[Unit]
Description="Trigger unconfining and re-confining user accounts on demand"
ConditionPathExists=/usr/bin/unconfinehomeusers

[Path]
PathChanged=/etc/passwd
PathChanged=/home
Unit=unconfinehomeusers.service

[Install]
WantedBy=multi-user.target
